# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
user=root
pidfile=/tmp/supervisord.pid

[program:uiserver]
command=uvicorn server:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
user=root
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:agent]
command=python agent.py start
directory=/app
user=root
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
